# config/systemd/kaizen.service
[Unit]
Description=Kaizen voice assistant
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/kaizen
ExecStart=%h/kaizen/run.sh --voice
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
